# profile.d/aliases.zsh
#!/bin/bash

for arg in crap shit masaneta cacete merda porra caralho cu cuuuuuuuu; do
  alias ${arg}="fuck"
done

alias ls="ls -G"
alias ll="ls -lah -G"
alias tailf='tail -f'
alias crontab="env EDITOR=vim crontab"

